# cfddns.timer
[Unit]
Description=Activate CFDDNS script daily

[Timer]
OnBootSec=1m
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target