# cfddns.service
[Unit]
Description=CloudFlare DDNS Script

[Service]
ExecStart=/usr/sbin/cfddns.py
Type=oneshot

[Install]
WantedBy=multi-user.target 